# ansible/roles/python-django/templates/coco-daphne.service.j2
[Unit]
Description=Coco TestAI ASGI Server (Daphne)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ app_user }}
Group={{ app_user }}
WorkingDirectory={{ backend_dir }}
Environment="HOME=/home/{{ app_user }}"
Environment="PATH={{ backend_dir }}/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="DJANGO_SETTINGS_MODULE={{ django_settings_module }}"
Environment="APP_ENV={{ env_name }}"
ExecStart={{ backend_dir }}/venv/bin/daphne -b 0.0.0.0 -p 8000 copilot_orchestrator.asgi:application
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
